# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 2 "app:create_app()"
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr 